# init/ams-publisher
#!/bin/bash
# ams-publisher daemon
# chkconfig: 345 20 80
# description: ams-publisher daemon
# processname: ams-publisherd  

. /etc/rc.d/init.d/functions
PROG_NAME="ams-publisherd"
DAEMON_PATH="/usr/bin/ams-publisherd"

case "$1" in
start)
	echo $"Starting $PROG_NAME: "
	msg=$(daemon $DAEMON_PATH -d start)
	retval=$?
	echo $msg
	exit $retval
;;
stop)
	echo $"Stopping $PROG_NAME: "
	msg=$(daemon $DAEMON_PATH -d stop)
	retval=$?
	echo $msg
	exit $retval
;;
restart)
	echo $"Restarting $PROG_NAME: "
	msg=$(daemon $DAEMON_PATH -d restart)
	retval=$?
	echo $msg
	exit $retval
;;
status)
	echo $"Status $PROG_NAME: "
	msg=$(daemon $DAEMON_PATH -d status)
	retval=$?
	echo $msg
	exit $retval
;;
*)
	echo "Usage: $0 {start|stop|restart|status}"
	exit 1
esac

